# speedtest.cron
# Run every 30 seconds
*/15 * * * * * * python3 ./main.py